# launchd/com.pipeline-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pipeline-bot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/ysh/dev/pipline-bot/venv/bin/python</string>
        <string>/Users/ysh/dev/pipline-bot/run.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/ysh/dev/pipline-bot</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>WEBHOOK_URL</key>
        <string>https://bot.ki-garu.com</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/Users/ysh/dev/pipline-bot/logs/launchd-bot.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ysh/dev/pipline-bot/logs/launchd-bot-error.log</string>
</dict>
</plist>
